# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/lib/proxy/supervisord.log

[program:proxy]
command=ncs-proxy run
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=redis-server
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0